# sbin/cleanup.crontab
00 6 * * * /bin/bash -l -c "%APP_HOME%/PROD/bin/cleanup.sh"
